# deploy/legendary-spoon.service
[Unit]
Description=legendary-spoon - Personal Toolkit Web Application
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/path/to/legendary-spoon
Environment="NODE_ENV=production"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/local/bin/bun run src/index.ts
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 安全选项
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/path/to/legendary-spoon/data /path/to/legendary-spoon/uploads

[Install]
WantedBy=multi-user.target
