# deploy/legendary-spoon.service
[Unit]
Description=legendary-spoon - Personal Toolkit Web Application
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/path/to/legendary-spoon
Environment="NODE_ENV=production"
ExecStart=/bin/bash /path/to/legendary-spoon/start.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 基本安全选项
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
